# scripts/distribution/debian/debian_openport/postinst
update-rc.d openport defaults
openport --restart-shares
