# scripts/distribution/redhat/fs/etc/init/openport.conf
# Openport - Starts the openport processes for the users in /etc/openport/users.conf

start on runlevel [2345]

task

console output 

exec /usr/bin/openport --restart-shares
