# scripts/debian/postinst
update-rc.d openport-manager defaults
/etc/init.d/openport-manager start
