# scripts/distribution/debian/debian/postinst
update-rc.d openport defaults
/etc/init.d/openport start
